# raspberry/var/spool/cron/crontabs/tc
*/1 * * * * ssh -f -N -R 60107:localhost:22 root@portal.it-port.ru
